# bmcweb.service.in
[Unit]
Description=Start bmcweb server

Wants=network.target
After=network.target
After=xyz.openbmc_project.User.Manager.service

[Service]
ExecReload=kill -s HUP $MAINPID
ExecStart=@MESON_INSTALL_PREFIX@/bin/bmcweb
Type=simple
WorkingDirectory=/home/root

[Install]
WantedBy=network.target
