# config/bmcweb.service.in
[Unit]
Description=Start bmcweb server

Wants=network.target
After=network.target
After=xyz.openbmc_project.User.Manager.service
After=xyz.openbmc_project.State.BMC.service
After=xyz.openbmc_project.Software.BMC.Updater.service
After=xyz.openbmc_project.State.Host@0.service
After=srvcfg-manager.service

[Service]
ExecReload=kill -s HUP $MAINPID
ExecStart=@MESON_INSTALL_PREFIX@/bin/bmcweb
Type=simple
WorkingDirectory=/home/root

[Install]
WantedBy=network.target
