# config/bmcweb.service.in
[Unit]
Description=Start bmcweb server

Wants=network.target
After=network.target

[Service]
ExecReload=kill -s HUP $MAINPID
ExecStart=@MESON_INSTALL_PREFIX@/bin/bmcweb daemon
Type=simple
User=@BMCWEB_USER@
Group=@BMCWEB_GROUP@
DynamicUser=@BMCWEB_DYNAMICUSER@
StateDirectory=@BMCWEB_STATE_DIRECTORY@
SyslogLevelPrefix=true
WatchdogSec=@BMCWEB_WATCHDOG_TIMEOUT_SECONDS@s
TemporaryFileSystem=/tmp/bmcweb:noexec,nodev,nosuid
WorkingDirectory=@BMCWEB_WORKING_DIRECTORY@

# bmcweb currently uses /tmp as a mechanism to share files.
PrivateTmp=no

[Install]
WantedBy=network.target
