# config/bmcweb.socket.in
[Unit]
Description=BMC Webserver socket

[Socket]
ListenStream=@BMCWEB_PORT@
ReusePort=true
Service=bmcweb.service
FileDescriptorName=bmcweb_@BMCWEB_PORT@_@HTTP_LEVEL_ALLOWED@_@HTTP_AUTH_LEVEL@
BindToDevice=@HTTP_BIND@

[Install]
WantedBy=sockets.target
