# config/bmcweb.socket.in
[Unit]
Description=BMC Webserver socket

[Socket]
ListenStream=@HTTPS_PORT@
ReusePort=true

[Install]
WantedBy=sockets.target
