# config/pam_auth.socket.in
[Unit]
Description=BMC Webserver authentication unix socket

[Socket]
ListenStream=@bmcweb_auth

[Install]
WantedBy=sockets.target
